# liveview.service
[Unit]
Description=Capture Movie Stream and Play with MPlayer
After=network.target

[Service]
Type=simple
ExecStart=/bin/sh -c 'gphoto2 --capture-movie --stdout | mplayer -cache 128 -'
User=pi
Restart=on-failure

[Install]
WantedBy=multi-user.target
